# .config/systemd/user/kanshi.service
[Unit]
Description=Start the kanshi profile service

[Service]
ExecStart=/usr/bin/kanshi

[Install]
WantedBy=default.target
